# SystemConfig/systemd/overrides/prometheus.service.d/cpu-affinity.conf
[Service]
CPUAffinity=2
MemoryHigh=1G
MemoryMax=2G
