# SystemConfig/systemd/services/configure-nic-irq-affinity.service
[Unit]
Description=NIC IRQ Affinity Configuration
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/configure-nic-irq-affinity.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
